# build/php/supervisor.conf
[program:sidekicka-app-def]
process_name=%(program_name)s_%(process_num)02d
numprocs=1
command=php /var/www/artisan queue:work redis --queue=default --timeout=1800
autostart=true
autorestart=true
redirect_stderr=true
user=root
stdout_logfile=/var/www/storage/logs/queue-default.log
stdout_logfile_maxbytes=10MB
logfile_backups=10
